# qhtlfirewall/qhtlwaterfall.service
[Unit]
Description=QhtLink Firewall - qhtlwaterfall
After=network.target qhtlfirewall.service

[Service]
LimitNOFILE=4096
ExecStart=/usr/sbin/qhtlwaterfall
ExecStartPre=/bin/rm -f /var/run/qhtlwaterfall.pid
Type=forking
PIDFile=/var/run/qhtlwaterfall.pid
KillMode=process
TimeoutStartSec=30s
TimeoutStopSec=30s
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
